# syncpass.service
[Unit]
Description=Sync passwords

[Service]
Type=simple
ExecStart=/home/jschwab/dotfiles/syncpass
