# synccal.timer
[Unit]
Description=Sync calendar files hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
