# syncorg.timer
[Unit]
Description=Sync org files hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
